# timescaledb/rootfs/etc/s6-overlay/s6-rc.d/pgagent/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start PgAgent service if enabled
# ==============================================================================
declare update
declare version_file

update=false;
version_file=/data/version

# Check for updates
if [[ $(< ${version_file}) != $(bashio::addon.version) ]]; then
	update=true
fi

# Wait for postgres to become available..
while ! psql -U "postgres" postgres -c "" 2> /dev/null; do
	sleep 1
done

# Check if we need to update..
if bashio::var.true "${update}"; then
	bashio::log.info "PgAgent: Waiting for update to complete.."
	while [[ $(< ${version_file}) != $(bashio::addon.version) ]]; do
		sleep 5
		bashio::log.info "PgAgent: Waiting for update to complete.."
	done

	bashio::log.info "Updating PgAgent Extension.."
	psql -U "postgres" -d "postgres" -X -c "ALTER EXTENSION pgagent UPDATE;" || true >> /var/log/pgagent.run.log
fi

if ! (psql -X -U postgres -d "postgres" -c "select 1 from pg_extension where extname = 'pgagent';" | grep -q 1); then
	bashio::log.info "Enabling PgAgent Extension.."
	psql -U "postgres" -d "postgres" -c "CREATE EXTENSION IF NOT EXISTS pgagent;" || true >> /var/log/pgagent.run.log
fi;

# Create a tmpfs for the pgagent log
mkdir -p "/data/tmp/pgagent"
export TMPDIR=/data/tmp/pgagent

# And run the daemon in the foreground
bashio::log.info "Starting PgAgent.."
/usr/local/bin/pgagent hostaddr=127.0.0.1 dbname=postgres user=postgres -f -l 1 #-s /var/log/pgagent
