# timescaledb/rootfs/etc/s6-overlay/s6-rc.d/init-addon/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: TimescaleDb
# Initializes the container during startup
# ==============================================================================
declare postgres_data
declare version_file
declare new_install

postgres_data=/data/postgres
version_file=/data/version
new_install=false;

# Applies permission to the data directory for the postgres user
applyPermissions () {
	chown -R postgres:postgres ${postgres_data}
	chmod 700 ${postgres_data}
}

# Initializes the data directory
initializeDataDirectory () {
	# Init data-directory
    bashio::log.info "Initializing new postgres directory.."
	mkdir -p ${postgres_data}
	applyPermissions
	su - postgres -c "initdb -D ${postgres_data}"
	# Set timescaledb as being enabled in the postgres config file.
	sed -r -i "s/[#]*\s*(shared_preload_libraries)\s*=\s*'(.*)'/\1 = 'timescaledb,\2'/;s/,'/'/" ${postgres_data}/postgresql.conf
	
	# Set Password protect IPv4 hosts by default
	echo "host    all             all             0.0.0.0/0               md5" >> ${postgres_data}/pg_hba.conf
	echo "local    all             all                                    md5" >> ${postgres_data}/pg_hba.conf
	echo "local    all             all                                   peer" >> ${postgres_data}/pg_hba.conf
	# Set Listen on all addresses (*)
	sed -r -i "s/[#]listen_addresses.=.'.*'/listen_addresses\ \=\ \'\*\'/g" ${postgres_data}/postgresql.conf
	# Set telemetry level
	echo "timescaledb.telemetry_level=$(bashio::config 'timescaledb.telemetry')" >> ${postgres_data}/postgresql.conf
	echo "timescaledb.license_key='CommunityLicense'" >> ${postgres_data}/postgresql.conf
	bashio::log.info "done"
}

# Upgrades the timescale extensions on all databases.
upgradeTimeScaleExtension () {
	# Upgrade Timescale..
	bashio::log.info "Upgrading Timescale extentions.."
	bashio::log.info "Updating Timescale Extension for system databases.."

	# Fetch all databases..
	psql \
		-X \
		-U "postgres" \
		-c "select datname from pg_database where datallowconn = true;" \
		--set ON_ERROR_STOP=off \
		--no-align \
		-t \
		--field-separator ' ' \
		--quiet \
	| while read datname; do
		psql -X -U postgres -d ${datname} -c "select 1 from pg_extension where extname = 'timescaledb';" | grep -q 1 \
		&& (
			bashio::log.info "Try updating Timescale Extension for database: '${datname}'.."; \
			# Drop the old functions and views, these could be incompatible with the new version.
			# See: https://github.com/timescale/timescaledb-extras/issues/34
			# Please note that we do not drop the extension itself, because that would drop the hypertables as well.
			# The views and functions could be recreated from: https://github.com/timescale/timescaledb-extras
			psql -U "postgres" homeassistant -c "drop function if exists public.get_dimension_details(regclass)"; \
			psql -U "postgres" homeassistant -c "drop view if exists public.chunks_tstz"; \
			psql -U "postgres" homeassistant -c "drop view if exists public.chunks_ts"; \
			(psql -X -U postgres -d ${datname} -c "ALTER EXTENSION timescaledb UPDATE;" || true)
		)
	done

	bashio::log.info "done"
}

# Upgrades the PostgreSQL databases
# $1: The current version
# $2: The new version
upgradePostgreSQL() {
	local fromVersion=${1}
	local toVersion=${2}

	bashio::log.notice "Upgrading databases now. This could take a while. Please be patient..."

	# Move the old data directory out of our way..
	mv ${postgres_data} ${postgres_data}${fromVersion}

	# And upgrade PostgreSQL
	bashio::log.notice "Upgrading PostgreSql from ${fromVersion} to ${toVersion}.."

	# Backup old HBA.conf and create a temp one...
	mv ${postgres_data}${fromVersion}/pg_hba.conf ${postgres_data}${fromVersion}/pg_hba_backup.conf
	echo "local    all             all                                     trust" > ${postgres_data}${fromVersion}/pg_hba.conf

	#Start postgres on the old data-dir
	bashio::log.info "Starting PostgreSQL-${fromVersion} first.."
	su - postgres -c "/usr/libexec/postgresql${fromVersion}/postgres -D ${postgres_data}${fromVersion}" &
	postgres_pid=$!

	# Wait for postgres to become available..
	while ! psql -X -U "postgres" postgres -c "" 2> /dev/null; do
		sleep 1
	done

	# Upgrade Timescale first, otherwise, pg_upgrade will fail.
	upgradeTimeScaleExtension

	# Stop server
	bashio::log.notice "Stopping PostgreSql ${fromVersion}.."
	kill ${postgres_pid}
	wait ${postgres_pid} || true
	bashio::log.notice "done"

	# Restore HBA.CONF
	rm ${postgres_data}${fromVersion}/pg_hba.conf
	mv ${postgres_data}${fromVersion}/pg_hba_backup.conf ${postgres_data}${fromVersion}/pg_hba.conf

	# Create a fresh data-directory
	initializeDataDirectory

	# And upgrade!
	bashio::log.notice "Upgrading databases.."
	cd ${postgres_data}${fromVersion}
	if su -c "pg_upgrade --old-bindir=/usr/libexec/postgresql${fromVersion} --new-bindir=/usr/libexec/postgresql${toVersion} --old-datadir=${postgres_data}${fromVersion} --new-datadir=${postgres_data} --link --username=postgres" -s /bin/sh postgres; then
		bashio::log.notice "PostgreSQL upgraded succesfully!"
		# Remove the old cluster...
		rm -r ${postgres_data}${fromVersion}
		return 0
	else
		# copy logs
		cp -r /data/postgres/pg_upgrade_output.d  ${postgres_data}${fromVersion} || true

		# Rollback..
		rm -r ${postgres_data}
		mv ${postgres_data}${fromVersion} ${postgres_data}

		bashio::log.error "PostgreSQL could not upgrade! Please inspect any errors in the lines above!"
		return 1
	fi
}


# Prepare /tmp/pg_stat_tmp directory
mkdir -p /tmp/pg_stat_tmp
chown postgres:postgres /tmp/pg_stat_tmp


if ! bashio::fs.directory_exists "${postgres_data}"; then
    bashio::log.info "Detected a fresh installation! Welcome! We're setting things up for you.."
    new_install=true
else
    touch ${version_file}
	# Always re-apply permissions, because they seem to be reset after a snapshot restore
	applyPermissions
fi

# Initialize for new installs
if bashio::var.true "${new_install}"; then
	touch /data/firstrun
	bashio::addon.version > ${version_file}
	initializeDataDirectory
else

	# Check if we need to restore again.
	if $(bashio::config.true 'retry_upgrade'); then
		if bashio::fs.directory_exists "${postgres_data}14"; then
			bashio::log.notice "An aborted upgrade from Postgres 14 was detected. Restoring.."
			rm -r ${postgres_data}
			mv ${postgres_data}14 ${postgres_data}
		else
			bashio::config.suggest.false 'retry_upgrade' 'This option is only for temporary reasons: to recover from a failed upgrade.'
		fi
	fi

	# Check if we need to upgrade from 12 to 15
	if [[ $(< ${postgres_data}/PG_VERSION) == "12" ]]; then
		bashio::log.notice "A database upgrade is required from Postgres 12."
		bashio::log.error "Upgrade was not succesful. This addon version does no longer support Postgres 12. Please upgrade to the previous version of the addon and upgrade to Postgres 14 first." 
		exit 1
	fi

	# Check if we need to upgrade from 14 to 15
	if [[ $(< ${postgres_data}/PG_VERSION) == "14" ]]; then
		bashio::log.notice "A database upgrade is required from Postgres 14."
		if upgradePostgreSQL 14 15; then
			bashio::log.notice "Upgrade was succesful"
		else
			bashio::log.error "Upgrade was not succesful."
			exit 1
		fi
	fi
fi

bashio::log.info "done"

# Apply TimescaleDb mem/cpu tuning settings
bashio::log.info "Tuning resources.."

# Always patch telemetry level
sed -r -i "s/timescaledb.telemetry_level.=.'.*'/timescaledb.telemetry_level=$(bashio::config 'timescaledb.telemetry')/g" ${postgres_data}/postgresql.conf

chmod 707 "/usr/share/timescaledb/002_timescaledb_tune.sh"
TS_TUNE_MEMORY=$(bashio::config 'timescaledb.maxmemory') \
	TS_TUNE_NUM_CPUS=$(bashio::config 'timescaledb.maxcpus') \
	POSTGRESQL_CONF_DIR=${postgres_data} \
	/usr/share/timescaledb/002_timescaledb_tune.sh
bashio::log.info "done"

# Appy max connections
bashio::log.info "Applying max connections.."
sed -i -e "/max_connections =/ s/= .*/= $(bashio::config 'max_connections')/" ${postgres_data}/postgresql.conf
bashio::log.info "done"
